# deploy/hermes-console.service
[Unit]
Description=Hermes Console
After=network.target hermes-gateway.service
Wants=hermes-gateway.service

[Service]
Type=simple
User=nick
WorkingDirectory=/home/nick/projects/hermes-console
EnvironmentFile=/home/nick/projects/hermes-console/.env
ExecStartPre=/home/nick/projects/hermes-console/deploy/verify-bundle.sh
ExecStart=/home/nick/.hermes/node/bin/node --import tsx server/production.ts
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
