# deploy/hermes-workspace.service
[Unit]
Description=Hermes Workspace
After=network.target hermes-gateway.service
Wants=hermes-gateway.service

[Service]
Type=simple
User=nick
WorkingDirectory=/home/nick/projects/hermes-workspace
EnvironmentFile=/home/nick/projects/hermes-workspace/.env
ExecStart=/home/nick/.hermes/node/bin/node --import tsx server/production.ts
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
